# files/reboot.service
[Unit]
Description=Reboot system
After=local-fs.target
Requires=local-fs.target

[Service]
Type=oneshot
ExecStart=/bin/systemctl reboot
